# carbon-base-macos/rootfs/etc/supervisor/conf.d/redis.conf
[program:redis]
command=/usr/local/bin/run-redis.sh
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/redis.out.log
stderr_logfile=/var/log/redis.err.log
stopsignal=TERM
